# templates/default/upstart/storm-ui.conf.erb
# This file is managed by Chef any changes made will be overwritten

author "Nicholas Verbeck <nerdynick@gmail.com>"
description "Starts up an instance of Storm's Web UI"

# What user and group to run as
setuid root
setgid root

#Make sure the network is online 1st
start on net-device-up

respawn
respawn limit 2 5

script
<%= render "env/storm.erb", :variables => { :user => @user, :storm_home => @storm_home, :java_lib_path => @java_lib_path } %>
    python $STORM_HOME/bin/storm ui
end script
