# templates/default/upstart/storm-nimbus.conf.erb
# This file is managed by Chef any changes made will be overwritten

author "Nicholas Verbeck <nerdynick@gmail.com>"
description "Starts up an instance of Storm's Nimbus Server"

# What user and group to run as
setuid root
setgid root

#Make sure the network is online 1st
start on net-device-up

respawn
respawn limit 2 5

pre-start script
<%= render "env/storm.erb", :variables => { :user => @user, :storm_home => @storm_home, :java_lib_path => @java_lib_path } %>
    # Clean up old local files. This prevents a problem with corrupt files.
    rm -Rf $STORM_HOME/storm-local/nimbus/*
end script

script
<%= render "env/storm.erb", :variables => { :user => @user, :storm_home => @storm_home, :java_lib_path => @java_lib_path } %>
    # General storm parameters
    python $STORM_HOME/bin/storm nimbus
end script
